# etc/systemd/system/sshd.service.d/override.conf
# Source: https://github.com/GrapheneOS/infrastructure/tree/95ed9b1befe12974f5ed13445f7786ed08bacfc0/etc/systemd/system/sshd.service.d

[Service]
LimitNOFILE=8192
ManagedOOMPreference=avoid