# etc/systemd/system/docker-auto-update@.service
[Unit]
Description=Update Docker continers defined in /srv/%I
Requires=network-online.target
Requisite=docker.service
After=network-online.target
After=docker.service

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/srv/%i
ExecStart=/usr/bin/docker image prune -f -a --filter 'until=240h'
ExecStart=/usr/bin/docker compose pull
ExecStart=/usr/bin/docker compose up -d